# public/init-conf/mount_create_couchdb.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Start couchDb"
author        "R&R"


start on starting boot-services

post-start script
  mount -o remount rw /
  DIRNAME=/var/run/couchdb
  mkdir -p $DIRNAME && chown -R couchdb:couchdb $DIRNAME
  /etc/rc.d/couchdb start &
  /sbin/iptables-restore /etc/rules
end script
